# make-deb/debian/postrm
#!/bin/sh

( systemctl stop psrtd > /dev/null 2>&1 || exit 0 )
( systemctl disable psrtd > /dev/null 2>&1 || exit 0 )
( userdel psrt > /dev/null 2>&1 || exit 0 )
